# ansible/roles/pwn/files/etc/profile.d/gef.sh
alias gdb-gef='/usr/bin/gdb -x /opt/tools/gef/gef.py -q '
