# ansible/roles/pwn/files/etc/profile.d/pwndbg.sh
alias gdb-pwndbg='/usr/bin/gdb -x /opt/tools/pwndbg/gdbinit.py -q '
alias pwndbg="gdb-pwndbg"
